# docker-cpp-service/resources/docker/s6-rc.d/app-service/finish
#!/bin/sh
tail -n 25 /var/log/service.log